# templates/supervisor.ini.j2
# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

[program:sentry-web]
autorestart=true
autostart=true
command={{sentry_home}}/sentry start
directory={{sentry_home}}
redirect_stderr=true
stdout_logfile={{sentry_home}}/logs/sentry-web.log
user={{sentry_user}}

[program:sentry-worker]
autorestart=true
autostart=true
command={{sentry_home}}/sentry celery worker -B
directory={{sentry_home}}
redirect_stderr=true
stdout_logfile={{sentry_home}}/logs/sentry-worker.log
user={{sentry_user}}
